# nats/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:nats]
command=bash -c "cd /home/nats && /usr/local/bin/bundle exec nats-server"