# gorouter/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:gorouter]
command=gorouter -c /home/gorouter/conf.yml